# keyboard/.config/systemd/user/k2.service
[Unit]
Description=Change keyboard layout automatically when pluging in keychron K2
 
[Service]
Type=simple
ExecStart=/home/jonathan/.local/bin/k2.py